# chromebook1/files/etc/cron.d/flyers
### flyers
SCRIPT_DIR=/home/sho/workspace/flyers
00 10,19 * * *  sho    cd "${SCRIPT_DIR}"; /home/sho/.local/share/virtualenvs/flyers-OhApUgB6/bin/activate; /usr/bin/python3 -m flyers >/dev/null 2>&1
