# chromebook1/files/etc/cron.d/hstn-family-tools
### garbage_notify
00 20 * * *       python3  /root/tools/hstn-family-scripts/garbage_notify.py >/dev/null 2>&1

### yurichan_okodukai
00 13 * * *       python3  /root/tools/hstn-family-scripts/yurichan_okodukai_remainder.py >/dev/null 2>&1

### chromebook1 server-stats
01 * * * *        python3  /root/tools/hstn-family-scripts/server-stats.py >/dev/null 2>&1
