# chromebook1/files/etc/cron.d/hstn-family-tools
### garbage_notify
00 20 * * *  sho    /usr/bin/python3  /home/sho/workspace/hstn-family-scripts/garbage_notify.py >/dev/null 2>&1
00 17 * * *  sho    /usr/bin/python3  /home/sho/workspace/hstn-family-scripts/garbage_notify.py >/dev/null 2>&1

### furikomi_reminder
00 13 * * *  sho    /usr/bin/python3  /home/sho/workspace/hstn-family-scripts/furikomi_remainder.py >/dev/null 2>&1

### chromebook1 server-stats
01 * * * *   sho    /usr/bin/python3  /home/sho/workspace/hstn-family-scripts/server-stats.py >/dev/null 2>&1
58 * * * *   sho    /usr/bin/python3  /home/sho/workspace/hstn-family-scripts/server-stats_speedtest.py >/dev/null 2>&1

### flayers
00 10,19 * * *  sho    /usr/bin/python3  /home/sho/workspace/hstn-family-scripts/flayers.py >/dev/null 2>&1
