# _archive/chromebook1/files/etc/systemd/system/netdata.service.d/override.conf
[Service]
ReadWriteDirectories=/var/spool/postfix/maildrop
